# starters/laravel/crontab
* * * * * cd /application && /usr/local/bin/php artisan schedule:run --no-interaction
